# profile.d/path.bat
set CMDER_USER_BIN="%~dp0..\%bin"
%lib_path% enhance_path_recursive "%CMDER_USER_BIN%" %max_depth%
