# profile.d/path.bat
set CMDER_CUSTOM_BIN="%~dp0..\%bin"
%lib_path% enhance_path_recursive "%CMDER_CUSTOM_BIN%" %max_depth%
